# roles/applications/templates/container-miniflux-app.service.j2
# container-miniflux-app.service
# autogenerated by Podman 3.1.2
# Mon May 31 14:28:59 BST 2021

[Unit]
Description=Podman container-miniflux-app.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target
RequiresMountsFor=/home/miniflux/.local/share/containers/storage /run/user/"{{ miniflux_user_uid }}"/containers
BindsTo=pod-miniflux.service
After=pod-miniflux.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStartPre=/bin/rm -f %t/container-miniflux-app.pid %t/container-miniflux-app.ctr-id
ExecStart=/usr/bin/podman container run --conmon-pidfile %t/container-miniflux-app.pid --cidfile %t/container-miniflux-app.ctr-id --cgroups=no-conmon --pod-id-file %t/pod-miniflux.pod-id --replace \
--name miniflux-app \
--restart=on-failure:5 \
--env BASE_URL=https://"{{ apps_domain }}"/ \
--env DATABASE_URL=postgres://miniflux:"{{ miniflux_postgres_password }}"@miniflux-postgres/miniflux?sslmode=disable \
--env RUN_MIGRATIONS=1 \
--env CREATE_ADMIN=1 \
--env ADMIN_USERNAME="{{ miniflux_admin_user }}" \
--env ADMIN_PASSWORD="{{ miniflux_admin_password }}" \
--expose 8080 \
--detach=True docker.io/miniflux/miniflux:latest
ExecStop=/usr/bin/podman container stop --ignore --cidfile %t/container-miniflux-app.ctr-id -t 10
ExecStopPost=/usr/bin/podman container rm --ignore -f --cidfile %t/container-miniflux-app.ctr-id
PIDFile=%t/container-miniflux-app.pid
Type=forking

[Install]
WantedBy=multi-user.target default.target
